# docker/systemd/es-inventory-ninja.timer
[Unit]
Description=Run ES Inventory Hub NinjaRMM Collector daily at 2 AM
Requires=es-inventory-ninja.service

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
