# temp_cross_vendor_timer.timer
[Unit]
Description=Run ES Inventory Hub Cross-Vendor Checks Daily at 3:00 AM
Requires=es-cross-vendor-checks.service

[Timer]
# Run daily at 3:00 AM Central Time (8:00 AM UTC)
OnCalendar=*-*-* 08:00:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
