# variances-hostnet.service
[Unit]
Description=Variances Dashboard (host-net)
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-/usr/bin/docker rm -f es-dashboards-variances-hostnet
ExecStart=/usr/bin/docker run -d --name es-dashboards-variances-hostnet --network host \
  --restart unless-stopped \
  -e FLASK_ENV=production \
  -e SECRET_KEY=variances-dashboard-secret-key \
  -e PORT=5001 \
  -e DB_DSN=postgresql://postgres:Xat162gT2Qsg4WDlO5r@localhost:5432/es_inventory_hub \
  -e THREATLOCKER_API_BASE_URL=https://portalapi.g.threatlocker.com \
  -e THREATLOCKER_API_KEY=24248BE80D8E0D4B85182FF90E7E141DBB99F90020F0B2C0D7C9269F1334CC24 \
  -e THREATLOCKER_ORGANIZATION_ID=dd850352-ee85-436b-8e41-818bdb52712c \
  variances-hostnet:test
ExecStop=/usr/bin/docker rm -f es-dashboards-variances-hostnet

[Install]
WantedBy=multi-user.target