# ops/systemd/es-inventory-api.service
[Unit]
Description=ES Inventory Hub API Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=rene
Group=rene
WorkingDirectory=/opt/es-inventory-hub
Environment=DB_DSN=postgresql+psycopg2://username:password@hostname:port/database_name
ExecStart=/usr/bin/python3 /opt/es-inventory-hub/api/api_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
