# ops/systemd/ninja-collector.timer
[Unit]
Description=Run Ninja collector daily at 02:30 AM Central Time

[Timer]
OnCalendar=*-*-* 02:30:00 America/Chicago
Persistent=true
AccuracySec=1m
Unit=ninja-collector.service

[Install]
WantedBy=timers.target
