# ops/systemd/sync-secrets.timer
[Unit]
Description=Timer for syncing API secrets from shared secrets
Documentation=https://github.com/es-inventory-hub/docs
Requires=sync-secrets.service

[Timer]
# Run daily at 3:00 AM Central Time (8:00 AM UTC)
OnCalendar=*-*-* 08:00:00
# Run immediately if missed
Persistent=true
# Add some randomization to avoid thundering herd
RandomizedDelaySec=300

[Install]
WantedBy=timers.target

