# ops/systemd/sync-secrets.service
[Unit]
Description=Sync API secrets from shared secrets to local .env
Documentation=https://github.com/es-inventory-hub/docs
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/es-inventory-hub
ExecStart=/usr/bin/python3 /opt/es-inventory-hub/scripts/sync_secrets.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sync-secrets

# Security settings
PrivateTmp=true
