# ops/systemd/windows-11-24h2-assessment.timer
[Unit]
Description=Run Windows 11 24H2 Assessment 45 minutes after Ninja collection

[Timer]
OnCalendar=*-*-* 23:45:00
Persistent=true

[Install]
WantedBy=timers.target
