# ops/systemd/windows-11-24h2-assessment.service
[Unit]
Description=Windows 11 24H2 Capability Assessment
# Run after any potential data collection

[Service]
Type=oneshot
User=postgres
Group=postgres
WorkingDirectory=/opt/es-inventory-hub
ExecStart=/opt/es-inventory-hub/.venv/bin/python3 /opt/es-inventory-hub/collectors/assessments/windows_11_24h2_assessment.py
StandardOutput=journal
StandardError=journal
Environment=DB_DSN=postgresql://postgres:mK2D282lRrs6bTpXWe7@localhost:5432/es_inventory_hub

[Install]
WantedBy=multi-user.target
