# systemd/es-secrets-reload.service
[Unit]
Description=Reload ES Inventory API after secrets change

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart es-inventory-api.service
