# systemd/es-secrets-monitor.path
[Unit]
Description=Monitor shared secrets for changes

[Path]
PathModified=/opt/shared-secrets/api-secrets.env
Unit=es-secrets-reload.service

[Install]
WantedBy=multi-user.target
